# .config/supervisord.conf
[program:nginx]
command=nginx -g "daemon off;"

[program:uwsgi]
command=uwsgi --ini /srv/tpay/.config/uwsgi.ini