# .config/supervisord.conf
[program:gunicorn]
command=gunicorn -c /srv/project/.config/gunicorn_cfg.py config.wsgi

[program:nginx]
command=nginx -g 'daemon off;'